# supervisor/ghost.conf
[program:ghost]
autorestart=true
directory=/usr/src/ghost
command=/entrypoint.sh npm start
